# deployer.service
[Unit]
Description=Deployer daemon
After=network.target

[Service]
User=deploy
Group=deploy
WorkingDirectory=/home/deploy
ExecStart=/home/deploy/deployer.exe --port 5000

[Install]
WantedBy=multi-user.target